# deploy/crontab
# Solar Directory Cron Jobs
# Managed by Docker — do not edit manually

# Flynax core cron (listing expiration, saved searches, cleanup)
*/30 * * * * www-data php /var/www/html/cron/index.php > /dev/null 2>&1

# Weekly pipeline: scrape 5 least-recent states (Sundays at 2 AM UTC)
0 2 * * 0 www-data /var/www/html/scripts/venv/bin/python /var/www/html/scripts/pipeline_orchestrator.py --mode weekly >> /var/www/html/scripts/logs/pipeline.log 2>&1

# Monthly verification: re-check all existing listings (1st of month at 3 AM UTC)
0 3 1 * * www-data /var/www/html/scripts/venv/bin/python /var/www/html/scripts/pipeline_orchestrator.py --mode monthly >> /var/www/html/scripts/logs/pipeline.log 2>&1
